# quadlets/opsechat-tor.container
# Quadlet file for Tor daemon container
# Place in ~/.config/containers/systemd/ or /etc/containers/systemd/

[Unit]
Description=Opsechat Tor Daemon Container
Documentation=https://github.com/HyperionGray/opsechat
After=network-online.target
Wants=network-online.target

[Container]
Image=alpine:latest
ContainerName=opsechat-tor
Exec=sh -c "apk add --no-cache tor netcat-openbsd && tor -f /etc/tor/torrc"
Volume=%h/opsechat/torrc:/etc/tor/torrc:ro,Z
Volume=opsechat-tor-data:/var/lib/tor:Z
Network=opsechat-network.network
HealthCmd=nc -z localhost 9051
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=10s
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target